# root/etc/services.d/sonarr/run
#!/command/with-contenv bash
umask "${UMASK}"
cd "${APP_DIR}/sonarr/bin" || exit
exec s6-setuidgid agpsn mono --debug "${APP_DIR}/sonarr/bin/Sonarr.exe" --nobrowser --data="${CONFIG_DIR}"

